# rpi/systemd/ptp4l.service.d/time-reference-monitor.conf
# Drop-in for ptp4l.service – Time Reference Monitor
#
# Overrides the ptp4l ExecStart to:
#   1. Use /etc/linuxptp/ptp4l.conf (free_running 1, slaveOnly 1 → monitor mode)
#   2. Read the PTP network interface from /etc/time-reference-monitor.conf
#
# PTP_IFACE must be set in /etc/time-reference-monitor.conf, e.g.:
#   PTP_IFACE=eth0
#
# The same interface must be set as --iface in time-reference-monitor.service.
#
# Installed to: /etc/systemd/system/ptp4l.service.d/time-reference-monitor.conf
# by setup.sh / update.sh
[Service]
EnvironmentFile=-/etc/time-reference-monitor.conf
ExecStart=
ExecStart=/usr/sbin/ptp4l -f /etc/linuxptp/ptp4l.conf -i ${PTP_IFACE:-eth0}
