# rpi/systemd/ptp4l.service.d/uds-permissions.conf
# Drop-in for ptp4l.service:
# Make both ptp4l Unix Domain Sockets world-accessible so the 'ptp' user
# can run pmc queries without root privileges.
#
#   /var/run/ptp4l    – management socket (uds_address, read/write)
#   /var/run/ptp4lro  – read-only socket  (uds_ro_address, used by this app)
#
# The monitor app connects to /var/run/ptp4lro via:
#   pmc -u -s /var/run/ptp4lro -b <domain> GET CURRENT_DATA_SET
#
# Installed to: /etc/systemd/system/ptp4l.service.d/uds-permissions.conf
# by setup.sh / update.sh
[Service]
ExecStartPost=/bin/sh -c 'chmod 666 /var/run/ptp4l 2>/dev/null; chmod 666 /var/run/ptp4lro 2>/dev/null; true'
