# rpi/systemd/ptp4l.service.d/uds-permissions.conf
# Drop-in for ptp4l.service:
# Make the ptp4l Unix Domain Socket world-accessible so the 'ptp' user
# can run pmc queries (GET DEFAULT_DATA_SET, etc.) without root privileges.
#
# Installed to: /etc/systemd/system/ptp4l.service.d/uds-permissions.conf
# by setup.sh / update.sh
[Service]
ExecStartPost=/bin/chmod 666 /var/run/ptp4l
